# scripts/netinsight-backend.service
[Unit]
Description=NetInsight Backend API Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/net-traffic/backend
Environment="PATH=/home/pi/net-traffic/backend/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/pi/net-traffic/backend/venv/bin/python3 /home/pi/net-traffic/backend/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings (commented out to allow packet capture)
# NoNewPrivileges=true
PrivateTmp=true

# Resource limits (adjust as needed)
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

